# roles/linx/templates/supervisor_confd_linx.conf.j2
[program:linx]
command={{ linx.path.exec }} -config {{ linx.path.config }}
directory={{ linx.path.home }}
